# supervisor/startup.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:php]
command=/etc/init.d/php7.0-fpm restart